# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:chainlit]
command=uv run chainlit run app.py --host 0.0.0.0 --port 8000
directory=/app
priority=1
autorestart=true
